# server/deploy/remote-controller-healthcheck.service
[Unit]
Description=Remote Controller healthcheck

[Service]
Type=oneshot
ExecStart=/opt/remote_controller/server/deploy/healthcheck.sh
