# server/deploy/systemd/remdesk-healthcheck.service
[Unit]
Description=RemDesk healthcheck

[Service]
Type=oneshot
ExecStart=/opt/remote_controller/server/deploy/scripts/healthcheck.sh
